# ldapauth/rootfs/etc/services.d/ldapauth/run
#!/usr/bin/with-contenv bash
exec python3 /usr/bin/ldap_server.py